# images/serf/supervisor/supervisor-serf.conf
[program:serf]
command     = serf-agent
numprocs    = 1
autostart   = true
autorestart = true